# May21/spring-boot-deployment/springpetclinic.service
[Unit]
Description=Spring Pet Clinic
After=syslog.target

[Service]
User=springbootuser
ExecStart=/usr/bin/java -jar /usr/share/springpetclinic/spring-petclinic-2.4.2.jar
SuccessExitStatus=143 

[Install] 
WantedBy=multi-user.target